# skf-os-labs-images/skf-openbox/root/etc/services.d/labs/run
#!/bin/execlineb -P

# Redirect stderr to stdout.
fdmove -c 2 1

# Wait until openbox is running
if { s6-svwait -t 10000 -U /var/run/s6/services/openbox/ }
# Wait until guacamole is running
if { s6-svwait -t 10000 -U /var/run/s6/services/guacamole/ }

# Drop privileges and set env
s6-setuidgid abc
s6-env DISPLAY=:1 HOME=/config

# Execute polybar
polybar -q main -c /config/.config/polybar/docky/config.ini 
