# skf-os-labs-images/root/etc/services.d/guacd/run
#!/bin/execlineb -P

# Redirect stderr to stdout.
fdmove -c 2 1

# Notify service manager when guacd is up
s6-notifyoncheck -w 500 -n 0 -c "redirfd -w 1 /dev/null fdmove -c 2 1 nc -w 1 -z -v 127.0.0.1 4822"

# Run guacd
/usr/sbin/guacd -f -b 127.0.0.1 -l 4822
